# templates/rc.local.erb
#!/bin/sh

### THIS FILE IS MANAGED BY PUPPET. ANY MANUAL CHANGES WILL GET OVERWRITTEN.
#
# Simplified "init" script for launching services, with very little tailoring
# needed for each OS.
#
# WARNING: If you want to customize /etc/rc.local (i.e. this file) just don't do
# it. Instead add you own scripts under /etc/rc.local.d and make them
# executable. This way they will be loaded automatically on boot.

SCRIPTDIR="/etc/rc.local.d"

# If no arguments (e.g. during startup), launch all services
if [ "$1" = "" ]; then
	run-parts $SCRIPTDIR

# If an argument is given, load only the appropriate script
else
	if [ -x $SCRIPTDIR/$1 ]; then
		$SCRIPTDIR/$1
	fi
fi

exit 0
